# Dailyquotes/crontab
0 7 * * * /app/exec.sh >> /app/logs/cron.log 2>&1
